# misc/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/10 * * * * /usr/bin/run-harbormaster > /proc/1/fd/1 2>/proc/1/fd/2
